# zram/rc.local.j2
#!/bin/bash
# chmod 0755 /etc/rc.local

/root/zram.sh + 600 0 {{ username }}

exit 0
